# deploy/systemd/charactier-backend-green.service
[Unit]
Description=Charactier AI Backend (Green)
After=network.target mongodb.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/charactier-ai/backend
EnvironmentFile=/var/www/charactier-ai/backend/.env
Environment="NODE_ENV=production"
Environment="PORT=5001"
Environment="PATH=/root/.nvm/versions/node/v22.16.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/root/.nvm/versions/node/v22.16.0/bin/node /var/www/charactier-ai/backend/dist/src/index.js
ExecReload=/bin/kill -USR2 $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=charactier-backend-green

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target