# deploy/systemd/charactier-backend-blue.service
[Unit]
Description=Charactier AI Backend (Blue)
After=network.target mongodb.service

[Service]
Type=simple
User=root
WorkingDirectory=/var/www/charactier-ai/backend
Environment="NODE_ENV=production"
Environment="PORT=5000"
Environment="PATH=/home/jun/.nvm/versions/node/v22.16.0/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/home/jun/.nvm/versions/node/v22.16.0/bin/node /var/www/charactier-ai/backend/dist/index.js
ExecReload=/bin/kill -USR2 $MAINPID
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=charactier-backend-blue

# Graceful shutdown
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target